# profile.d/etcd.sh
#!/bin/sh

export PATH=/usr/local/etcd:$PATH
